# desktop/autostart/plank.desktop
[Desktop Entry]
Name=plank
Comment=plank dock
Exec=plank
Type=Application
Terminal=false
Hidden=false

